# .config/systemd/user/alert-battery.service
[Unit]

[Service]
Type=oneshot
ExecStart=/home/jakob/scripts/alert-battery.sh

[Install]
WantedBy=graphical.target
